# systemd/.config/systemd/user/newsboat.service
[Unit]
Description=Update newsboat feed

[Service]
Type=simple
ExecStart=/usr/bin/newsboat -x reload
