# systemd/.config/systemd/user/thunderbird.timer
[Unit]
Description=Periodically update thunderbird for new emails

[Timer]
OnBootSec=10s
OnUnitActiveSec=30m

[Install]
WantedBy=timers.target
